# jenqins.service
[Unit]
Description=jenqins
After=network.target

[Service]
ExecStart=/home/username/jenqins/build/jenqins
